# docker/supervisord.conf
; =============================================================================
;  Jumpcoin Qt wallet - supervisord configuration
;  Runs as root; individual program blocks drop to the appropriate user
;  via the `user=` directive (jumpcoin for GUI/VNC, debian-tor for tor).
;
;  Process order (priority ascending, higher number = starts later):
;    10 Xvfb          - virtual X display
;    20 openbox       - lightweight window manager
;    25 pcmanfm       - desktop icon layer
;    30 x0vncserver   - TigerVNC attached to the running display
;    40 noVNC + tor   - browser VNC bridge / Tor daemon (tor variant only)
;    50 jumpcoin-qt   - the Qt GUI wallet (full node via server=1)
; =============================================================================

[supervisord]
nodaemon=true
# Run supervisord itself as root so it can setuid to any account
# (notably debian-tor, uid 101, for the bundled tor process in the
# :tor image variant).  Individual [program:*] blocks drop to the
# jumpcoin user as needed via `user=jumpcoin`.
user=root
logfile=/home/jumpcoin/.supervisor/supervisord.log
pidfile=/home/jumpcoin/.supervisor/supervisord.pid
loglevel=info
directory=/home/jumpcoin
environment=
    HOME="/home/jumpcoin",
    DISPLAY=":99",
    USER="jumpcoin"

[unix_http_server]
file=/home/jumpcoin/.supervisor/supervisor.sock
chmod=0700
chown=jumpcoin:jumpcoin

[supervisorctl]
serverurl=unix:///home/jumpcoin/.supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

; -----------------------------------------------------------------------------
;  Xvfb - virtual framebuffer (display :99)
; -----------------------------------------------------------------------------
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_DISPLAY_WIDTH)sx%(ENV_DISPLAY_HEIGHT)sx%(ENV_DISPLAY_DEPTH)s
    -nolisten tcp -dpi 96 +extension GLX +extension RANDR +extension RENDER
autostart=true
autorestart=true
priority=10
stopasgroup=true
killasgroup=true
environment=HOME="/home/jumpcoin"
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

; -----------------------------------------------------------------------------
;  Openbox - lightweight window manager
; -----------------------------------------------------------------------------
[program:openbox]
command=/usr/bin/openbox
autostart=true
autorestart=true
priority=20
stopasgroup=true
killasgroup=true
environment=DISPLAY=":99",HOME="/home/jumpcoin",USER="jumpcoin"
user=jumpcoin
stdout_logfile=/var/log/supervisor/openbox.log
stderr_logfile=/var/log/supervisor/openbox.err
startsecs=1

; -----------------------------------------------------------------------------
;  pcmanfm - draws clickable desktop icons (so the jumpcoin-qt.desktop
;  shortcut in ~/Desktop is actually usable).  Runs in --desktop mode
;  with no main window, only the wallpaper + icon layer.
; -----------------------------------------------------------------------------
[program:pcmanfm]
command=/usr/bin/pcmanfm --desktop
autostart=true
autorestart=true
priority=25
stopasgroup=true
killasgroup=true
environment=DISPLAY=":99",HOME="/home/jumpcoin",USER="jumpcoin",XDG_CURRENT_DESKTOP=XFCE
user=jumpcoin
stdout_logfile=/var/log/supervisor/pcmanfm.log
stderr_logfile=/var/log/supervisor/pcmanfm.err
startsecs=2

; -----------------------------------------------------------------------------
;  x11vnc - VNC server attached to the running Xvfb display.
;  x11vnc keeps X11 CLIPBOARD/PRIMARY selections in sync with VNC cut text,
;  which noVNC's clipboard panel relies on.
; -----------------------------------------------------------------------------
[program:vncserver]
command=/usr/local/bin/start-vnc
autostart=true
autorestart=true
startretries=10
startsecs=2
priority=30
stopasgroup=true
killasgroup=true
environment=
    HOME="/home/jumpcoin",
    DISPLAY=":99",
    VNC_PORT="%(ENV_VNC_PORT)s",
    VNC_SECURITY_TYPES="%(ENV_VNC_SECURITY_TYPES)s",
    VNC_INSECURE_OK="%(ENV_VNC_INSECURE_OK)s"
user=jumpcoin
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc.err

; -----------------------------------------------------------------------------
;  noVNC + websockify - browser-accessible VNC
; -----------------------------------------------------------------------------
[program:autocutsel]
command=/usr/bin/autocutsel -selection CLIPBOARD
autostart=true
autorestart=true
priority=35
stopasgroup=true
killasgroup=true
environment=DISPLAY=":99",HOME="/home/jumpcoin",USER="jumpcoin"
user=jumpcoin
stdout_logfile=/var/log/supervisor/autocutsel.log
stderr_logfile=/var/log/supervisor/autocutsel.err
startsecs=1

[program:novnc]
command=/opt/novnc/utils/novnc_proxy
    --listen %(ENV_NOVNC_PORT)s
    --vnc localhost:%(ENV_VNC_PORT)s
    --web /opt/novnc
autostart=true
autorestart=true
startretries=5
startsecs=2
priority=40
stopasgroup=true
killasgroup=true
environment=
    HOME="/home/jumpcoin",
    PYTHONPATH="/opt/websockify:/opt/novnc/utils"
user=jumpcoin
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err

; -----------------------------------------------------------------------------
;  jumpcoin-qt - the Qt GUI wallet (acts as both daemon and wallet)
;
;  Why not a separate jumpcoind?  Both processes open an exclusive
;  flock on the data directory (.lock), so running them side by side
;  always results in "Cannot obtain a lock on data directory".  Qt has
;  the full node built in (server=1 in jumpcoin.conf enables RPC) so
;  the standalone daemon is redundant here.
; -----------------------------------------------------------------------------
[program:jumpcoin-qt]
command=/usr/local/bin/start-jumpcoin-qt
autostart=true
autorestart=true
startsecs=5
startretries=20
priority=50
stopasgroup=true
killasgroup=true
environment=DISPLAY=":99",HOME="/home/jumpcoin",USER="jumpcoin"
user=jumpcoin
stdout_logfile=/var/log/supervisor/jumpcoin-qt.log
stderr_logfile=/var/log/supervisor/jumpcoin-qt.err

; -----------------------------------------------------------------------------
;  tor (only present in the :tor image variant)
;
;  We always declare the program block so the file works in both
;  variants.  The command is wrapped in a shell test that no-ops when
;  /usr/bin/tor is missing (i.e. the :latest build), so the block is
;  a no-op there.  In the :tor build /usr/bin/tor is installed and
;  the daemon starts before jumpcoin-qt (priority 40 < 50).
;
;  Supervisord runs as root, so it can setuid to debian-tor directly.
;  We removed the "User" directive from /etc/tor/torrc so tor does not
;  attempt its own privilege drop (which would fail inside a container).
; -----------------------------------------------------------------------------
[program:tor]
command=/bin/sh -c 'if [ -x /usr/bin/tor ]; then exec /usr/bin/tor -f /etc/tor/torrc; else exit 0; fi'
autostart=true
autorestart=unexpected
exitcodes=0
startsecs=3
startretries=10
priority=40
stopasgroup=true
killasgroup=true
user=debian-tor
stdout_logfile=/var/log/supervisor/tor.log
stderr_logfile=/var/log/supervisor/tor.err
